# ansible/playbooks/apps/jellyfin/mnt-media.automount.j2
[Unit]
Description=Auto-mount media USB drive

[Automount]
Where=/mnt/media

[Install]
WantedBy=multi-user.target
