# ansible/playbooks/resources/jellyfin/mnt-media.mount.j2
[Unit]
Description=Mount media USB drive

[Mount]
What=UUID=344EEF6D4EEF25F8
Where=/mnt/media
Type=ntfs
Options=noatime
