# ansible/playbooks/resources/jellyfin/mnt-media.mount.j2
[Unit]
Description=Mount media USB

[Mount]
What=/dev/disk/by-uuid/{{ media_disk_id }}
Where=/mnt/media
Type=ext4
Options=noatime,defaults

[Install]
WantedBy=multi-user.target
